# tg-canvas.service
# tg-canvas.service — copy to /etc/systemd/system/ and edit the placeholders
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now tg-canvas
#
# Replace:
#   YOUR_USER        → the system user that should run the server (e.g. wdai, ubuntu, pi)
#   YOUR_INSTALL_DIR → absolute path to the tg-canvas-miniapp directory

[Unit]
Description=OpenClaw Telegram Canvas Mini App
After=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=YOUR_INSTALL_DIR
EnvironmentFile=YOUR_INSTALL_DIR/.env
ExecStart=/usr/bin/node YOUR_INSTALL_DIR/server.js
Restart=always
RestartSec=2
StandardOutput=append:YOUR_INSTALL_DIR/logs/server.log
StandardError=append:YOUR_INSTALL_DIR/logs/server.log

[Install]
WantedBy=default.target
